# home/.config/autostart/compton.desktop
[Desktop Entry]
Exec=/bin/compton
Name=compton
Type=Application
