# home/.config/autostart/poly-centre.desktop
[Desktop Entry]
Exec=/bin/polybar centre
Name=Polybar centre
Type=Application
