# home/.config/autostart/poly-right.desktop
[Desktop Entry]
Exec=/bin/polybar right
Name=Polybar right
Type=Application
