# home/.config/autostart/poly-left.desktop
[Desktop Entry]
Exec=/bin/polybar left
Name=Polybar left
Type=Application
